# npoed-plp-edx-sysadminguide/src/plp-task-import-analytics.timer
[Unit]
Description=Import analytics

[Timer]
OnCalendar=*-*-* 5:00

[Install]
WantedBy=timers.target
